# roles/python-src/files/pythonapp/flask-gunicorn/supervisor/app.conf
[program:flaskapp]
command=/usr/bin/runapp.sh
user=www-data
autostart=true
autorestart=true
